# insights/backend/server/scripts/systemd/codeinsight-sync-commits.timer
[Unit]
Description=Schedule CodeInsight Stage A metadata sync

[Timer]
OnBootSec=1m
OnUnitActiveSec=10m
AccuracySec=1m
Persistent=true
Unit=codeinsight-sync-commits.service

[Install]
WantedBy=timers.target
